# rhcsa-labs.service
[Unit]
Description=RHCSA Practice Labs
After=network.target

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/rhcsa-practice-labs
Environment="PATH=/home/exedev/rhcsa-practice-labs/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="BIND=0.0.0.0:8080"
Environment="LOG_LEVEL=info"
ExecStart=/home/exedev/rhcsa-practice-labs/.venv/bin/gunicorn \
    --config gunicorn.conf.py \
    api.app_socketio:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
